# sfplib/rootfs/etc/s6-overlay/s6-rc.d/init-sfplib/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: SFPLiberate
# Runs initialization commands before services start
# ==============================================================================
set -e

bashio::log.info "Initializing SFPLiberate..."

# Create necessary directories
bashio::log.info "Creating data directories..."
mkdir -p /config/sfplib/submissions
mkdir -p /config/sfplib/backups

bashio::log.info "Initialization complete"
